# misc/ansible-config-manager.service
[Unit]
Description=Ansible Config Manager

[Service]
ExecStart=/usr/local/bin/ansible-config-manager

[Install]
WantedBy=default.target
